# templates/units/%name%.target
[Unit]
Description={name} services

[Install]
WantedBy=default.target

